# view-oj.ini
[program: view-oj-app]
command = /home/view-oj/venv/bin/gunicorn -w 8 -b 0.0.0.0:5000 flask_app:app
directory = /home/view-oj
autostart = true
redirect_stderr = true
stdout_logfile = /etc/supervisor/log/view-oj-app.log

[program: view-oj-celery]
command = /home/view-oj/venv/bin/celery -A tasks worker -l info -c 8 --pool=eventlet
directory = /home/view-oj
autostart = true
redirect_stderr = true
stdout_logfile = /etc/supervisor/log/view-oj-celery.log

[program: view-oj-beat]
command = /home/view-oj/venv/bin/celery -A tasks beat -l info
directory = /home/view-oj
autostart = true
redirect_stderr = true
stdout_logfile = /etc/supervisor/log/view-oj-beat.log